# arm/alpine-arm/deluge-alpine-arm/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=info

[program:deluged]
command=deluged -d

[program:deluge-web]
command=deluge-web

[program:tor]
command=tor
user=tor
